# scheduler.crontab
0 18 * * * /home/fsnowdin/Music/organize.sh >> /home/fsnowdin/Music/cronlog.txt
